# deploy/pharos-ki-proxy.service
[Unit]
Description=Pharos KI-Proxy (Claude Vision/Enrich + Übergabe-Trainer)
After=network.target

[Service]
Type=simple
User=deploy
WorkingDirectory=/opt/pharos
# --env-file lädt ANTHROPIC_API_KEY etc. vor dotenv/config (server/ki-proxy.mjs)
ExecStart=/usr/bin/node --env-file=/opt/pharos/.env.local server/ki-proxy.mjs
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
